# delete_requests.cron
00 * * * * curl -X DELETE "http://${DISPATCHER_IP}:${DISPATCHER_PORT}/api/delete"